# deploy/systemd/geomaker-touchterrain-watchdog.timer
[Unit]
Description=Timer para TouchTerrain Watchdog (a cada 2 minutos) — SPEC-935-R145

[Timer]
OnBootSec=1min
OnUnitActiveSec=2min
AccuracySec=10s

[Install]
WantedBy=timers.target
