# deploy/systemd/geomaker-touchterrain-watchdog.service
# Referência/documentação — o instalador (wsl/setup.sh) gera esta unidade
# automaticamente. Não é necessário editar ou instalar este arquivo manualmente.
[Unit]
Description=TouchTerrain Watchdog — verifica /main e reinicia geomaker-touchterrain se travado (SPEC-935-R145)
After=geomaker-touchterrain.service

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /opt/geomaker/installer/touchterrain_watchdog.py --url http://localhost:8081/main --service geomaker-touchterrain --timeout 10
User=root
